# legacy/pattern-strategy-fastapi/trading-bot/tools/launchd/com.tradingbot.newsagent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.tradingbot.newsagent</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>cd "__BOT_DIR__" &amp;&amp; exec python3 news_agent.py</string>
  </array>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>15</integer>

  <key>StandardOutPath</key>
  <string>__BOT_DIR__/state/launchd_newsagent.log</string>
  <key>StandardErrorPath</key>
  <string>__BOT_DIR__/state/launchd_newsagent.log</string>
</dict>
</plist>
